# nnic-request.service
[Unit]
Description=initiates the request module of the recorder.

[Service]
User=root
ExecStart=/usr/bin/php /opt/nnic/request.php
